# ansible_conf/templates/mathematic_site.ini
[program:mathematic_site]
command=/usr/local/bin/uwsgi
  --socket 127.0.0.1:9000
  --protocol uwsgi
  --python-path /home/www-mathematic/mathematic_site
  --virtualenv /home/www-mathematic/mathematic_site_venv
  --wsgi-file website/__init__.py
  --callable app
  --need-app
  --master
  --processes 1
  --threads 10
  --log-x-forwarded-for
directory=/home/www-mathematic/mathematic_site
autostart=true
autorestart=true
user=www-mathematic
